# rootfs/etc/s6-overlay/s6-rc.d/svc-nicotine/run
#!/command/with-contenv bash

export HOME=/data

# Start Nicotine
exec s6-setuidgid $(id -nu $PUID) nicotine